# deploy/eodhd-download.service
# User-level service — install to ~/.config/systemd/user/
# See deploy/eodhd-download.timer for the schedule.
[Unit]
Description=Turtle EODHD Data Download
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/turtle-backtest
EnvironmentFile=%h/.config/turtle-backtest/secrets.env
ExecStart=%h/.local/bin/uv run python scripts/download_eodhd_data.py --data history
StandardOutput=journal
StandardError=journal
