# deploy/snapshot_company.timer
# Monthly company snapshot timer — user-level systemd timer (no sudo required)
#
# DEPLOY
#   mkdir -p ~/.config/systemd/user
#   ln -s ~/turtle-backtest/deploy/snapshot_company.service ~/.config/systemd/user/snapshot_company.service
#   ln -s ~/turtle-backtest/deploy/snapshot_company.timer   ~/.config/systemd/user/snapshot_company.timer
#   systemctl --user daemon-reload
#   systemctl --user enable --now snapshot_company.timer
#
# PREREQUISITES
#   secrets file must exist at ~/.config/turtle-backtest/secrets.env (see Phase 3 of docs/implementation.md)
#   loginctl enable-linger <username>   # allows timer to fire without an active login session
#
# STATUS / LOGS
#   systemctl --user status snapshot_company.timer
#   systemctl --user list-timers snapshot_company.timer
#   journalctl --user -u snapshot_company.service
#
# REMOVE
#   systemctl --user disable --now snapshot_company.timer
#   rm ~/.config/systemd/user/snapshot_company.{service,timer}
#   systemctl --user daemon-reload

[Unit]
Description=Monthly turtle.company snapshot

[Timer]
OnCalendar=*-*-01 01:00:00
TimeZone=Europe/Tallinn
Persistent=true

[Install]
WantedBy=timers.target
